# mcp_center/servers/oe_cli_mcp_server/mcp-server.service
[Unit]
Description=MCP Tool Registration Service
After=network.target
After=multi-user.target
# 新增：确保文件系统就绪（避免工作目录未挂载）
RequiresMountsFor=/usr/lib/euler-copilot-framework/mcp_center/

[Service]
User=root
Group=root
# 工作目录：必须是项目根目录（已正确配置，确保 server.py 中相对路径生效）
WorkingDirectory=/usr/lib/euler-copilot-framework/mcp_center/

Environment="PATH=/usr/lib/euler-copilot-framework/mcp_center/servers/oe_cli_mcp_server/venv/global/bin:$PATH"

ExecStart=/usr/lib/euler-copilot-framework/mcp_center/servers/oe_cli_mcp_server/venv/global/bin/python /usr/lib/euler-copilot-framework/mcp_center/servers/oe_cli_mcp_server/server.py

# 原有合理配置保留
Restart=always
RestartSec=5
KillMode=control-group
Environment="LANGUAGE=zh"
Environment="LOG_LEVEL=INFO"
StandardOutput=journal+console
StandardError=journal+console


[Install]
WantedBy=multi-user.target