# mcp_center/service/nload-server.service
[Unit]
Description=Nload Command Server for MCP Center
After=network.target

[Service]
User=root
WorkingDirectory=/usr/lib/euler-copilot-framework/mcp_center
# 关键：添加 PYTHONPATH，将工作目录加入模块搜索路径
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/usr/lib/euler-copilot-framework/mcp_center
ExecStart=/usr/bin/python3 servers/nload/src/server.py
Restart=always
RestartSec=5
KillMode=control-group   
KillSignal=SIGTERM      
TimeoutStopSec=30    
SuccessExitStatus=143 

[Install]
WantedBy=multi-user.target
