# witty_mcp_manager/data/witty-mcp-manager.service
[Unit]
Description=Witty MCP Manager - Universal MCP Host/Loader
Documentation=https://atomgit.com/openeuler/euler-copilot-framework
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=witty-mcp
Group=witty-mcp
WorkingDirectory=/usr/lib/witty-mcp-manager
Environment=HOME=/var/lib/witty-mcp-manager
Environment=WITTY_MCP_CONFIG=/etc/witty/mcp-manager.yaml
Environment=UV_PROJECT_ENVIRONMENT=/var/lib/witty-mcp-manager/.venv
Environment=UV_CACHE_DIR=/var/lib/witty-mcp-manager/cache/uv
ExecStartPre=/usr/bin/bash -lc 'if [ -f /usr/lib/witty-mcp-manager/uv.lock ]; then exec /usr/bin/uv sync --locked --no-dev --no-editable --project /usr/lib/witty-mcp-manager; else exec /usr/bin/uv sync --no-dev --no-editable --project /usr/lib/witty-mcp-manager; fi'
ExecStart=/usr/bin/bash -lc 'if [ -f /usr/lib/witty-mcp-manager/uv.lock ]; then exec /usr/bin/uv run --locked --no-sync --project /usr/lib/witty-mcp-manager witty-mcp daemon; else exec /usr/bin/uv run --no-sync --project /usr/lib/witty-mcp-manager witty-mcp daemon; fi'
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/usr/lib/witty-mcp-manager
ReadWritePaths=/run /var/lib/witty-mcp-manager

# Resource limits
LimitNOFILE=65536
TasksMax=512

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=witty-mcp-manager

[Install]
WantedBy=multi-user.target
